# etc/systemd/system/limbo-backup.timer
[Unit]
Description=Run limbo-backup tasks daily

[Timer]
OnCalendar=daily
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
